# build/image/supervisord.conf
[supervisord]
nodaemon=true

[program:worker-service]
command=/usr/src/app/build/bin/main
directory=/usr/src/app
autorestart=true
autostart=true

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
